# s6.git-sync/etc/services.d/git-sync/run
#!/usr/bin/with-contenv bash
cd /
exec \
        s6-setuidgid abc /opt/bin/ssh-sync.sh
