# lazylibrarian.cron
#!/bin/bash
set -e
DOCKER_REPO=linuxserver/lazylibrarian

# Get CURRENT
CURRENT=$(curl https://hub.docker.com/v2/repositories/${DOCKER_REPO}/tags/latest | jq -c --raw-output '.last_updated')
if [[ -n "${CURRENT}" && "${CURRENT}" != "null" && "${CURRENT}" == "2"* ]]; then
	echo "${CURRENT}" > ${CONTAINER}.LATEST
	echo "${CURRENT}" >> ${CONTAINER}.VERSIONS
fi
