# s6.lakkris-old/etc/services.d/serf/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/serf agent --config-dir=/etc/serf
