# s6.shinobi/etc/services.d/shinobi-cron/run
#!/bin/execlineb -P

if { s6-svwait -U /var/run/s6/services/shinobi/ }

s6-setuidgid abc

/usr/bin/node /opt/shinobi/cron.js
