# s6.shinobi/etc/services.d/shinobi-motion/run
#!/bin/execlineb -P

if { s6-svwait -U /var/run/s6/services/shinobi/ }

cd /opt/shinobi

s6-setuidgid abc

/usr/bin/node /opt/shinobi/plugins/motion/shinobi-motion.js
