# lakkris/etc/services.d/inotifyd-core/run
#!/usr/bin/with-contenv bash
inotifyd /lakkris/.notify/router.sh /tmp/lakkris.notify:c
