# Dockerfiles/shinobi/etc/services.d/shinobi/run
#!/usr/bin/with-contenv bash
cd /opt/shinobi
exec \
	s6-setuidgid abc /usr/bin/node /opt/shinobi/camera.js
