# lakkris/etc/services.d/inotifyd/run
#!/usr/bin/with-contenv bash
inotifyd /lakkris/handlers/inotifyd.sh $(cat /tmp/lakkris.notify | cut -d',' -f-1)
